# board/zotac/overlay/etc/rc.local
#!/bin/sh -e
# CRDA workaround
export COUNTRY=US
/sbin/crda

# mount debugfs
mount -t debugfs none /sys/kernel/debug

# update links and cache to the
# available shared libraries
ldconfig

exit 0
